# cpp_file_server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log

[program:fileservice]
command=./server
